# roles/linux_common/files/add_alb_hosts_local.service
[Unit]
Description=Add ALB domains to /etc/hosts

[Service]
Type=oneshot
ExecStart=/usr/local/bin/add_alb_hosts_local.sh